# loops/com.ghq.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.ghq.scheduler.plist -- launchd agent for GHQ autonomous scheduler

  Runs scheduler.sh every 15 minutes via launchd on macOS.

  IMPORTANT: Before installing, replace YOUR_GHQ_PATH below with your actual
  GHQ repo path (e.g., /Users/you/repos/ghq).

  Installation:
    cp loops/com.ghq.scheduler.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.ghq.scheduler.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.ghq.scheduler.plist
    rm ~/Library/LaunchAgents/com.ghq.scheduler.plist

  Check status:
    launchctl list | grep com.ghq.scheduler

  Manual trigger:
    launchctl start com.ghq.scheduler

  Logs:
    tail -f /tmp/ghq-scheduler.log
    tail -f /tmp/ghq-scheduler-err.log
-->
<plist version="1.0">
<dict>
    <!-- Unique label for this launch agent -->
    <key>Label</key>
    <string>com.ghq.scheduler</string>

    <!-- Script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>export PATH="/usr/local/bin:/opt/homebrew/bin:$PATH"; exec "${GHQ_ROOT}/loops/scripts/scheduler.sh" 2>&amp;1</string>
    </array>

    <!-- Run every 15 minutes (900 seconds) -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>GHQ_ROOT</key>
        <string>YOUR_GHQ_PATH</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/ghq-scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ghq-scheduler-err.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>YOUR_GHQ_PATH</string>

    <!-- Keep alive: restart if it crashes, but not if it exits cleanly -->
    <key>KeepAlive</key>
    <false/>

    <!-- Don't run at load time (wait for first interval) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Timeout: kill if running longer than 10 minutes -->
    <key>TimeOut</key>
    <integer>600</integer>

    <!-- Low priority to avoid interfering with user work -->
    <key>LowPriorityIO</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
